# systemc/etc/systemc/system/getty@tty1.service
[Service]
ExecStart=bin/agetty --noclear
