# systemc/etc/systemc/system/rednet.service
[Unit]
Wants=graphical.target

[Service]
ExecStart=bin/texec rednet os.loadAPI("/rom/apis/rednet") rednet.run()
